# Bert/configuration/src/main/bin/bert-server
#!/bin/sh
# This is the init script for startup of the "bert" application.
# Logs from this script will be found in /var/log/bert/
#
# This script starts the server for "bert". This is necessary to
# operate either the bluetooth command or interactive terminal interfaces.

### BEGIN INIT INFO
# Provides:          bert
# Required-Start:    $local_fs $syslog
# Required-Stop:     $local_fs $syslog
# Should-Start:      $named
# Should-Stop:       $named
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start the bert control application
### END INIT INFO

set -e
export BERT_HOME=/usr/local/robot
export BERTLOGDIR=/var/log/bert
export BERTLOG=${BERTLOGDIR}/bert.log
export DAEMON=${BERT_HOME}/bin/bert
export PARAMS=""
PATH=/sbin:/bin:/usr/sbin:/usr/bin:${BERT_HOME}/bin

. /lib/lsb/init-functions
mkdir -p ${BERT_HOME}/logs
mkdir -p ${BERTLOGDIR}

case "$1" in
  start)
    # Uncomment for debugging, execute as: sudo bin/bert start
    #  set -x
	  rm -rf ${BERTLOG}
	  > ${BERTLOG}

	  # This is important. We need to wait until the IP Address of the machine
	  # is actually configured. Loop until the IP comes up. This is needed for
	  # the Bluetooth connection.
	  while true; do
	    IP="`hostname -I`"
	    if [ "$IP" ] ; then
	      echo "IP = ${IP}" >> ${BERTLOG}
	      break
	    fi
	    sleep 1
	  done

    IP=`echo ${IP}|sed -e 's/^[ \t]*//'`
	  echo "===================== /etc/init.d/bert-server start ===================" >> ${BERTLOG}
    # start-stop-daemon will only start bert if there is no running instance
	  start-stop-daemon --start --background $DAEMON
	  log_end_msg $?
  ;;

  restart)
    /etc/init.d/bert-server stop
    /etc/init.d/bert-server start
  ;;

  stop)
	  start-stop-daemon --stop --retry 5 $DAEMON
	  log_end_msg $?
	  echo "=================== /etc/init.d/bert-server stopped ===============" >> ${BERTLOG}
  ;;
  force-reload)
    /etc/init.d/bert-server restart
  ;;

  *)
    echo "Usage: /etc/init.d/bert-server {start|stop|restart|force-reload}"
    exit 1
    ;;
esac
